# INSTALL/supervisor_config/timeline_piper.conf
#
# /etc/supervisor/conf.d/timeline_piper.conf
#
# Generated by Chef for ip-172-27-85-169.vpc.odesk.com
#

[program:timeline_piper]
command=/usr/bin/php /var/www/html/timeline/event_logs/log_monitor_piper.php
autostart=true
autorestart=true
startretries=3
user=root
